# systemd/backup.service
[Unit]
Description=Backup to ftp

[Service]
ExecStart=/usr/local/bin/backup.sh
Type=simple