# udev/80-powersave.rules
ACTION=="change",SUBSYSTEM=="power_supply",RUN+="/usr/bin/kill -ALRM \"$(cat \"/var/run/dwmbar/bar.pid")\""
