# udev/80-powersave.rules
ACTION=="change",SUBSYSTEM=="power_supply",RUN+="/usr/bin/kill -USR1 \"$(cat \"/tmp/dwmbar.lock\")\""
