# udev/80-powersave.rules
ACTION=="change",SUBSYSTEM=="power_supply",RUN+="/usr/bin/pkill -USR1 dwmbar.sh"
